# etc/init.d/apparmor-openrc
#!/sbin/openrc-run

description="Loads AppArmor security profiles"
APPARMOR_FUNCTIONS="/lib/apparmor/rc.apparmor.functions"

depend() {
	need localmount procfs
	after modules
	before local
	provide apparmor
}

aa_log_daemon_msg() {
	einfo "$@"
}

aa_log_action_start() {
	ebegin "$@"
}

aa_log_action_end() {
	eend "$@"
}

aa_log_success_msg() {
	eend 0 "$@"
}

aa_log_failure_msg() {
	eend 1 "$@"
}

aa_log_warning_msg() {
	ewarn "$@"
}

_is_unprivileged_container() {
	if [ -f /.dockerenv ] || grep -q -E '/docker|/lxc' /proc/1/cgroup; then
		if ! is_container_with_internal_policy; then
			return 0
		fi
	fi
	return 1
}

start() {
	if [ ! -f "${APPARMOR_FUNCTIONS}" ]; then
		eerror "AppArmor function library not found at ${APPARMOR_FUNCTIONS}"
		return 1
	fi
	# shellcheck source=/dev/null
	. "${APPARMOR_FUNCTIONS}"

	if _is_unprivileged_container; then
		einfo "Not starting AppArmor in this container."
		return 0
	fi
	
	apparmor_start
}

stop() {
	einfo "Leaving AppArmor profiles loaded."
	ewarn "To tear down all profiles, run 'aa-teardown' manually."
	return 0
}

reload() {
	if [ ! -f "${APPARMOR_FUNCTIONS}" ]; then
		eerror "AppArmor function library not found at ${APPARMOR_FUNCTIONS}"
		return 1
	fi
	# shellcheck source=/dev/null
	. "${APPARMOR_FUNCTIONS}"

	if _is_unprivileged_container; then
		einfo "Not reloading AppArmor in this container."
		return 0
	fi

	apparmor_restart
}

status() {
	if [ ! -f "${APPARMOR_FUNCTIONS}" ]; then
		eerror "AppArmor function library not found at ${APPARMOR_FUNCTIONS}"
		return 1
	fi
	# shellcheck source=/dev/null
	. "${APPARMOR_FUNCTIONS}"

	apparmor_status
}
w